# scripts/systemd/claude-usage-report.service
[Unit]
Description=Claude usage weekly report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Optional: override GH_TOKEN for gh CLI authentication
EnvironmentFile=-%h/.config/claude-usage-report/env
ExecStart=%h/.local/bin/claude-usage-report
Restart=on-failure
StandardOutput=journal
StandardError=journal
