# deploy/systemd/claude-usage-tracker.timer
[Unit]
Description=Claude usage snapshot timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
