# deploy/launchd/com.user.claude-usage-tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-usage-tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.local/bin/claude-usage-tracker-snapshot</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_USAGE_TRACKER_DB_PATH</key>
        <string>{{HOME}}/.local/share/claude-usage-tracker/snapshots.db</string>
        <!-- Run `claude-usage-tracker-setup` to get correct values for your account -->
        <key>CLAUDE_USAGE_TRACKER_WEEKLY_RESET_DAY</key>
        <string>Thursday</string>
        <key>CLAUDE_USAGE_TRACKER_WEEKLY_RESET_HOUR</key>
        <string>7</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.local/share/claude-usage-tracker/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/.local/share/claude-usage-tracker/launchd.error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
